# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:postgres]
command=/bin/bash -c 'exec /usr/lib/postgresql/*/bin/postgres -D /data/postgres -k /var/run/postgresql'
user=postgres
priority=10
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:authentik-server]
command=/lifecycle/ak server
user=authentik
directory=/
priority=30
autorestart=true
startretries=30
environment=HOME="/data/authentik"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; The server applies DB migrations on startup; hold the worker back until
; they exist (capped at ~5 min, after which supervisord's restarts take over).
; A stale pidfile from a crashed run would make the worker refuse to start,
; so clear it first.
[program:authentik-worker]
command=/bin/bash -c 'for i in $(seq 60); do psql -h /var/run/postgresql -U authentik -d authentik -tAc "SELECT 1 FROM authentik_tenants_tenant LIMIT 1" >/dev/null 2>&1 && break; sleep 5; done; rm -f "${TMPDIR:-/tmp}/authentik-worker.pid"; exec /lifecycle/ak worker'
user=authentik
directory=/
priority=40
autorestart=true
startretries=30
environment=HOME="/data/authentik"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; Same start command as the official Grist image's sandbox/run.sh
[program:grist]
command=/usr/local/bin/node _build/stubs/app/server/server.js
directory=/grist
user=grist
priority=30
autorestart=true
startretries=30
environment=HOME="/data/grist",NODE_PATH="_build:_build/ext:_build/stubs:ext/node_modules"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"
priority=50
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
